# config/common/systemd/redirect.service
[Unit]
Description=Redirect

[Service]
ExecStart=/var/www/redirect/current/bin/redirect
Restart=always

[Install]
WantedBy=multi-user.target