# config/common/systemd/redirect.api.service
[Unit]
Description=Redirect

[Service]
ExecStart=/var/www/redirect/current/bin/api
Restart=always
User=redirect

[Install]
WantedBy=multi-user.target
