# config/common/systemd/redirect.api.service
[Unit]
Description=Redirect

[Service]
ExecStart=/var/www/redirect/current/bin/api --config-file /var/www/redirect/config.cfg --secret-file /var/www/redirect/secret.cfg --mail-dir /var/www/redirect/current/emails
Restart=always
User=redirect

[Install]
WantedBy=multi-user.target
